# common/pvr-env/autobuild/overrides/etc/profile.d/pvr.sh
export KWIN_COMPOSE=O2ES
export QT_XCB_GL_INTEGRATION=xcb_egl
