# common/pvr-env/autobuild/overrides/etc/profile.d/pvr.sh
export __EGL_VENDOR_LIBRARY_FILENAMES=/usr/share/glvnd/egl_vendor.d/50_pvr.json
export KWIN_COMPOSE=O2ES
export QT_QUICK_BACKEND=software
